# seagull/root/etc/services.d/seagull/run
#!/usr/bin/with-contenv sh

exec /seagull;
