# elasticsearch/root/etc/services.d/elasticsearch/run
#!/usr/bin/env bash

exec /elasticsearch/bin/elasticsearch;
