# elasticsearch/root/etc/services.d/elasticsearch/run
#!/usr/bin/env bash

exec s6/elasticsearch/bin/elasticsearch;
